# server-configs/mcp-meet.service
[Unit]
Description=Google Meet MCP Server
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/workspace_mcp
Environment=PYTHONPATH=/home/ubuntu/workspace_mcp
EnvironmentFile=/home/ubuntu/workspace_mcp/.env
ExecStart=/home/ubuntu/workspace_mcp/venv/bin/python /home/ubuntu/workspace_mcp/mcp_server/run_meet_mcp.py
Restart=always

[Install]
WantedBy=multi-user.target
